# templates/systemd/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
ListenStream=9443
